# hyprland/kwallet.service
[Unit]
Description=KDE Wallet Servcie

[Service]
ExecStart=/usr/bin/kwalletd5

[Install]
WantedBy=default.target
